# etc/crontabs/root
*/5 * * * * /etc/lpyscripts/checkfrpc.sh
